# .config/autostart/xmodmap.desktop
[Desktop Entry]
Name=xmodmap
Exec=bash -c "xmodmap $HOME/.Xmodmap"
Terminal=true
Type=Application
