# .config/autostart/xmodmap.desktop
[Desktop Entry]
Name=xmodmap
Exec=sh -c "xmodmap $HOME/.Xmodmap"
Terminal=true
Type=Application
